# pi-printer-client/systemd/byteorder-ble-printer.service
[Unit]
Description=ByteOrder BLE Print Server
After=bluetooth.service
Wants=bluetooth.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ble-print-server
Environment=PRINT_HOST=127.0.0.1
ExecStart=/opt/ble-print-server/venv/bin/python /opt/ble-print-server/print_server.py --port 8080
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=byteorder-ble-printer

[Install]
WantedBy=multi-user.target
